# meta-ref-os-iot/meta-5xx/recipes-bsp/5xx-udev-rules/files/hdmi-hotplug@.service
[Unit]
Description=Handle HDMI hotplug events

[Service]
Type=simple
ExecStart=/usr/bin/hdmi-hotplug.sh %I
